# systemd/openclaw-guardian.service
[Unit]
Description=OpenClaw Guardian - Health Monitor and Auto-Recovery
Documentation=https://github.com/yourusername/openclaw-guardian
After=network.target
# Note: OpenClaw gateway is managed by 'openclaw gateway start' command
# not by this service. We just need network to be up.

[Service]
Type=oneshot
User=%I
Environment="HOME=/home/%I"
Environment="PATH=/home/%I/.nvm/versions/node/v25.4.0/bin:/usr/local/bin:/usr/bin:/bin:/home/%I/.local/bin"

# Run the guardian check
ExecStart=/home/%I/Documents/GitHub/openclaw-guardian/bin/guardian.sh check

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-guardian

[Install]
WantedBy=multi-user.target
